# Linux system开发/build rootfs/ubuntu1604/src/components/dependencies/openvpn.conf
[program:openvpn]
command=openvpn --config client.ovpn
directory=/etc/openvpn
user=root
autorestart=true
